# deploy/systemd/ops-scheduler.service
[Unit]
Description=ops-agent scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ops-agent
Group=ops-agent
WorkingDirectory=/opt/ops-agent
EnvironmentFile=/etc/ops-agent/ops-agent.env
ExecStart=/usr/local/bin/ops-scheduler --env ${OPS_SCHEDULER_ENV} --env-file /etc/ops-agent/environments.yaml --policy /etc/ops-agent/policies.yaml --audit /var/lib/ops-agent/scheduler.db --audit-driver sqlite --incident-state-file /var/lib/ops-agent/incidents.db --notify-config /etc/ops-agent/notifications.yaml --notify-trigger-after 2 --notify-recovery-after 2 --discover-interval 6h
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/ops-agent

[Install]
WantedBy=multi-user.target
